# nginx/lavaland.service
[Unit]
Description=WSGI gunicorn
After=network.target

[Service]
User=lawrence
Group=www-data
WorkingDirectory=/sites/lavaland/python
ExecStart=/home/lawrence/.local/bin/gunicorn --workers 3 --bind unix:lavaland.sock -m 777 wsgi:app
Restart=always
KillSignal=SIGQUIT
Type=simple
StandardError=syslog
NotifyAccess=all

[Install]
WantedBy=multi-user.target
